# services/netns-nft@.service.d/bring-up-lo.conf
[Unit]
Requires=netns-lo@%i.service
After=netns-lo@%i.service
