# services/netns-lo@.service.d/part-of-macvlan.conf
[Unit]
# Bring lo up and down together with macvlan.
PartOf=netns-macvlan@%i.service
